# dotfiles/etc/systemd/system/miner.service
[Unit]
Description=Run on Boot

[Service]
Type=forking
ExecStart=/home/travis/launcher.sh

[Install]
WantedBy=multi-user.target
