# dotfiles/etc/systemd/system/software-update.service
[Unit]
Description=Software updates

[Service]
ExecStart=/home/travis/.bin/software-update

[Install]
WantedBy=multi-user.target
