# dotfiles/etc/udev/rules.d/10-local.rules
# Lock computer on any hardrive additions
#KERNEL=="sd[a-z]", RUN+="/home/travis/Documents/development/n0v1c3/linux/scripts/bash/i3-usb_lock.sh"
#KERNEL=="usb*", RUN+="/home/travis/Documents/development/n0v1c3/linux/scripts/bash/i3-usb_lock.sh"
